# roles/avans/templates/student_project_backup.service.j2
[Unit]
Description=Backup
Wants=student_project_backup.timer

[Service]
Type=oneshot
ExecStart=/bin/bash /etc/systemd/user/student_project_backup.sh

[Install]
WantedBy=multi-user.target

