# roles/avans/templates/student_project_backup.timer.j2
[Unit]
Description=Backup
Requires=student_project_backup.service

[Timer]
Unit=student_project_backup.service
OnCalendar=sat 12:00
Persistent=true

[Install]
WantedBy=timers.target