# roles/Cocalc_backup/templates/Cocalc_backup.timer.j2
[Unit]
Description=Cocalc_backup
Requires=Cocalc_backup.service

[Timer]
Unit=Cocalc_backup.service
OnCalendar=sun 12:00
Persistent=true

[Install]
WantedBy=timers.target